# other/pokedialog.service
[Unit]
Description=pokedialog service
After=network.target
StartLimitIntervalSec=0
[Service]
Environment=BOT_TOKEN=TOBEDEFINED
Type=simple
Restart=always
RestartSec=30
WorkingDirectory=/home/pi/pokedialog
User=pi
ExecStart=/home/pi/pokedialog/pokedialogbot

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pokedialog

StartLimitBurst=5

TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
